# container/systemd/logrotate/logrotate-custom.timer
[Unit]
Description=Run logrotate every week

[Timer]
OnBootSec=1min
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target