# container/systemd/logrotate/logrotate-custom.service
[Unit]
Description=Run custom logrotate

[Service]
Type=oneshot
ExecStart=/usr/sbin/logrotate /etc/logrotate.d/metric_log.conf